# docs/ops/systemd-telegram.service
[Unit]
Description=Codex Autorunner Telegram bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/you/car-workspace
EnvironmentFile=%h/.config/codex-autorunner/codex-autorunner.env
# Replace the ExecStart path with `which car` or the codex-autorunner binary.
ExecStart=/home/you/.local/bin/car telegram start --path /home/you/car-workspace
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
